# momentum.conf
[program:momentum]
command=momentum-prod --no-sandbox
autorestart=true